# huwai/supervisord.conf
[program:mytick]
command = /home/dotcloud/current/tick
